# IMSDK/check.crontab
*/1 * * * * /opt/IMSDK/check.sh  >/dev/null 
